# djangoproject/supervisor.ini
[program:brewpi]
command=/brew-ve/bin/gunicorn_django -c /brew-project/djangoproject/gunicorn_conf.py
directory=/brew-project/djangoproject/
user=vagrant
autostart=true
autorestart=true
priority=991
stopsignal=KILL